# supervisord.conf
[supervisord]
nodaemon=true

[program:webssh]
dictory = /var/www/app/
command =python /var/www/app/main.py


[program:nginx]
command=nginx
